# files/etc/uci-defaults/99-init-settings
#!/bin/sh
cat << 'EOF' > /root/.profile
export SHELL=/usr/bin/zsh
exec /usr/bin/zsh
EOF
chmod 644 /root/.profile
uci set luci.main.mediaurlbase='/luci-static/argon'
uci commit luci

exit 0